# launchd/com.shiftly.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shiftly.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/osascript</string>
    <string>/CHANGE_ME/Shiftly/scripts/sync.applescript</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SHIFTLY_ROOT</key>
    <string>/CHANGE_ME/Shiftly</string>
    <!-- Legacy names kept while scripts still read them. -->
    <key>SHIFTY_ROOT</key>
    <string>/CHANGE_ME/Shiftly</string>
    <key>SHIFTFLOW_ROOT</key>
    <string>/CHANGE_ME/Shiftly</string>
    <key>SHIFTLY_SYNC_MODE</key>
    <string>next_month</string>
    <key>SHIFTY_SYNC_MODE</key>
    <string>next_month</string>
    <key>SHIFTFLOW_SYNC_MODE</key>
    <string>next_month</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>28</integer>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/CHANGE_ME/Shiftly/data/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/CHANGE_ME/Shiftly/data/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
