# launchd/com.shiftly.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shiftly.sync</string>

  <!-- Headless sync through the app binary (same EventKit engine as the
       GUI). Grant Calendar access by syncing once from the app before
       loading this job. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Shiftly.app/Contents/MacOS/Shiftly</string>
    <string>--sync</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SHIFTLY_ROOT</key>
    <string>/CHANGE_ME/ShiftlyData</string>
    <key>SHIFTLY_SYNC_MODE</key>
    <string>next_month</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>28</integer>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/com.shiftly.sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.shiftly.sync.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
